# linux/autostart/ignore-lid-switch-tweak.desktop
[Desktop Entry]
Type=Application
Name=ignore-lid-switch-tweak
Exec=/usr/lib/gnome-tweak-tool/gnome-tweak-tool-lid-inhibitor 

