# linux/autostart/keybase_autostart.desktop
# This file is autogenerated when Keybase starts. Keybase will recreate it if you
# delete it, so to disable autostart you should comment out the Exec line below.
[Desktop Entry]
Name=Keybase autostart
Type=Application
NoDisplay=true
Exec=run_keybase
